# backend/shared-infrastructure/knowledge-base-services/kb-7-terminology/docker/kb7-supervisor.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:redis]
command=redis-server /etc/redis/redis-kb7.conf
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/redis-kb7.log
stderr_logfile=/var/log/redis-kb7.log
priority=100

[program:graphdb]
command=/opt/graphdb/bin/graphdb -Djava.awt.headless=true -Xmx3g -Xms1g
directory=/opt/graphdb
autostart=true
autorestart=true
stdout_logfile=/var/log/graphdb.log
stderr_logfile=/var/log/graphdb.log
priority=200
startsecs=30

[program:sparql-proxy]
command=/app/sparql-proxy/sparql-proxy
directory=/app/sparql-proxy
autostart=true
autorestart=true
stdout_logfile=/var/log/sparql-proxy.log
stderr_logfile=/var/log/sparql-proxy.log
priority=300
environment=GRAPHDB_URL="http://localhost:7200",REDIS_URL="localhost:6379",PORT="8095"

[program:kb7-init]
command=/app/init-kb7.sh
directory=/app
autostart=true
autorestart=false
exitcodes=0
stdout_logfile=/var/log/kb7-init.log
stderr_logfile=/var/log/kb7-init.log
priority=400
startsecs=0